# easystroke.service
[Unit]
Description=Easystroke
Documentation=man:easystroke

[Service]
Type=simple
ExecStart=/usr/bin/easystroke

[Install]
WantedBy=xinitrc.target
