# getmail.service
[Unit]
Description=Get mail
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=/usr/bin/getmail
