# xinitrc.target
[Unit]
Description=Xinitrc stuff
Wants=wm.service
