# xinitrc.target
[Unit]
Description=Xinitrc stuff
Requires=wm.service
