# xmonad.service
[Unit]
Description=Xmonad window manager

[Service]
ExecStart=/usr/bin/xmonad
 
[Install]
Alias=wm.service
